# docs/blog/posts/2025/012-consul-tf-backend/20-bootstrap-consul/roles/acme/files/acme-backend.service
[Unit]
Description=ACME backend
After=network.target

[Service]
Type=simple
ExecStartPre=/bin/mkdir -p /var/www/html
ExecStart=/usr/bin/python3 -m http.server 9000 -d /var/www/html

[Install]
WantedBy=multi-user.target
